# agent/deploy/lemma-agent.service.tmpl
# Lemma Agent — systemd unit.
#
# Generated by `lemma agent install --shape systemd --output <dir>`.
# Install with:
#
#   sudo cp lemma-agent.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now lemma-agent
#
# Adjust ExecStart paths and environment to match your deployment.
[Unit]
Description=Lemma Agent — federated compliance evidence agent
Documentation=https://github.com/JoshDoesIT/Lemma
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{BINARY_PATH}} serve --port={{HEALTH_PORT}} --evidence-dir={{EVIDENCE_DIR}} --keys-dir={{KEYS_DIR}}
Restart=on-failure
RestartSec=5

# Run as a dedicated, non-privileged user; create with:
#   sudo useradd --system --no-create-home --shell /usr/sbin/nologin lemma-agent
User=lemma-agent
Group=lemma-agent

# Hardening — the agent only needs read access to the keys dir,
# read+write to the evidence dir, and outbound network for `forward`.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ReadWritePaths={{EVIDENCE_DIR}}
ReadOnlyPaths={{KEYS_DIR}}

[Install]
WantedBy=multi-user.target
